# todoapp.socket
[Unit]
Description=My To-Do App

[Socket]
ListenStream=80